# examples/systemd/spotifycation.timer
[Unit]
Description=spotification timer

[Timer]
OnUnitActiveSec=30seconds

[Install]
WantedBy=timers.target

