# examples/systemd/spotifycation.timer
[Unit]
Description=spotification timer

[Timer]
OnBootSec=30seconds
OnUnitActiveSec=30seconds

[Install]
WantedBy=timers.target

